# template/master/etc/systemd/system/flanneld.service.d/50-network-config.conf.template
[Unit]
Requires=etcd-member.service
[Service]
ExecStartPre=/usr/bin/etcdctl set /coreos.com/network/config '{"Network":${POD_NETWORK}, "Backend": {"Type": "vxlan"}}'
